# conf/supervisord-ftp.conf
[program:vsftpd]
command=/usr/sbin/vsftpd /etc/vsftpd.conf
numprocs=1
autostart=true
autorestart=true

